# deploy/macos/com.inferra.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.inferra.agent</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-lc</string>
      <string>/usr/local/lib/inferra/inferra --config /usr/local/etc/inferra/inferra.toml init-db &amp;&amp; exec /usr/local/lib/inferra/inferra --config /usr/local/etc/inferra/inferra.toml serve</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/usr/local/var/log/inferra.log</string>
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/inferra.err.log</string>
  </dict>
</plist>
